# aaaa.conf
[program:aaaa]
command=python3 /home/curator2/aaaa/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/aaaa/test.err.log
stdout_logfile=/home/curator2/aaaa/test.out.log